# debian/data/etc/init/grails-blog.conf
description "Blog front end"
author "Charlie Knudsen <charlie.knudsen@gmail.com>"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

script
  # All application configuration done with environment variables
  export GITHUB_CLIENT_ID="xxxxxxxxxxxxxxx"
  export GITHUB_SECRET="xxxxxxxxxxxxxxxxx"
  export CLIENT_PREFIX="http://127.0.0.1:5678"
  export ETCD_URL="http://10.0.2.15:4001"
  export HOST_IP=127.0.0.1

  # command line values
  export GC_OPTS="-XX:+UseParallelOldGC"
  export JAVA_OPTS="-server -Xms256M -Xmx512M"
  export JETTY_JAR="/usr/share/jetty-runner-9.1.0/jetty.jar"
  export APP_WAR="/usr/share/grails-blog/application.war"
  export JAVA_HOME=/usr/lib/jvm/java-7-openjdk-amd64/jre

  $JAVA_HOME/bin/java $JAVA_OPTS $GC_OPTS -jar $JETTY_JAR $APP_WAR
end script
