# debian/data/etc/init/grails-blog.conf
description "Blog front end"
author "Charlie Knudsen <charlie.knudsen@gmail.com>"

start on (filesystem and net-device-up)
stop on runlevel [!2345]

respawn
respawn limit 5 30

script
  export GITHUB_CLIENT_ID="xxxxxxxxxxxxxxx"
  export GITHUB_SECRET="xxxxxxxxxxxxxxxxx"
  export CLIENT_PREFIX="http://127.0.0.1:5678"
  export ETCD_URL="http://10.0.2.15:4001"
  export HOST_IP=`/usr/bin/iplookup 'eth0 '`

  export GC_OPTS="-XX:+UseParallelOldGC"
  export JAVA_OPTS="-server -Xms256M -Xmx512M"
  export JETTY_JAR="/opt/jetty-runner/jetty-runner-9.1.0.jar"
  export APP_WAR="/opt/grails-blog/blog-0.1-SNAPSHOT.war"

  /usr/bin/java $JAVA_OPTS $GC_OPTS -jar $JETTY_JAR $APP_WAR
end script